# docker/s6-overlay/s6-rc.d/setup/up
#!/bin/execlineb -P
/bin/sh /etc/s6-overlay/s6-rc.d/setup/init.sh
