# docker/s6-overlay/s6-rc.d/log-sender/run
#!/bin/sh
# Wait for config to be present before starting
while [ ! -f /etc/log-sender/config.json ]; do
  sleep 5
done
echo "Starting log-sender service..." >> /data/logs/startup.log
exec gosu nonroot log-sender service --config-file /etc/log-sender/config.json >> /data/logs/log-sender.log 2>&1
