# docker/s6-overlay/s6-rc.d/logrotate-cron/run
#!/bin/sh
while true; do
  logrotate /etc/logrotate.d/holochain.conf
  sleep 86400
done
